# contrib/debian/preremove.sh
#!/bin/bash
systemctl stop saos-agent
systemctl disable saos-agent
